# supervisor.conf
[supervisord]
nodaemon=true

[program:sugar]
directory=/node_app
command=/node_app/node_modules/.bin/coffee /node_app/index.coffee
redirect_stderr=true
autostart=true
autorestart=true
